# centos6-supervisord/files/skel.conf
[program:progname]
command = command arg1 arg2
directory = /path/to/workdir
user = progusername
autostart = false
autorestart = true