# docker-lab/dockerfiles/shinken/snmpd.conf
[program:snmpd]
command = /usr/sbin/snmpd
autorestart = true
log_stdout=true
log_stderr=true
