# shinken-1.4.4/files/supervisor/broker.conf
[program:broker]
command=/usr/bin/python /usr/bin/shinken-broker -r -c /etc/shinken/brokerd.ini
autostart=true
user=shinken
log_stdout=true
log_stderr=true
