# shinken-2.0.3-webui-dev/files/supervisor/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --dbpath /var/lib/mongo
autostart=true
user=mongod
log_stdout=true
log_stderr=true
