# shinken-2.0.3-webui-dev/files/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
log_stdout=true
log_stderr=true